# conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nginx]
command=service apache2 start
autorestart=false
redirect_stderr=true
